# projects/network/vpn-manager/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:dockerd]
command=dockerd --host=unix:///var/run/docker.sock --storage-driver=overlay2
stdout_logfile=/var/log/supervisor/dockerd.log
stderr_logfile=/var/log/supervisor/dockerd.log
autorestart=true
priority=100

[program:vpn-manager]
command=/app/vpn-manager.sh
stdout_logfile=/var/log/supervisor/vpn-manager.log
stderr_logfile=/var/log/supervisor/vpn-manager.log
autorestart=true
priority=200
depends_on=dockerd 