# configs/systemd/rclone-zurg-mount.service
[Unit]
Description=RClone Mount for Zurg Docker Data
Documentation=https://rclone.org/docs/
After=network-online.target
Wants=network-online.target
AssertPathIsDirectory=/mnt/rclone

[Service]
Type=notify
User=root
Group=root
ExecStartPre=/bin/mkdir -p /mnt/rclone/docker
ExecStart=/usr/bin/rclone mount zurg:docker_data /mnt/rclone/docker \
    --allow-non-empty \
    --allow-other \
    --umask=002 \
    --attr-timeout 8760h \
    --buffer-size 64M \
    --dir-cache-time 120s \
    --poll-interval 60s \
    --vfs-cache-max-age 2m \
    --vfs-cache-max-size 30G \
    --vfs-cache-min-free-space 1G \
    --vfs-cache-mode full \
    --vfs-cache-poll-interval 30s \
    --vfs-disk-space-total-size 32G \
    --vfs-fast-fingerprint \
    --vfs-read-ahead 64M \
    --vfs-read-chunk-size 1M \
    --vfs-read-chunk-size-limit 32M \
    --vfs-read-wait 40ms \
    --vfs-refresh \
    --transfers 16 \
    --checkers 16 \
    --multi-thread-streams 8
ExecStop=/bin/fusermount -u /mnt/rclone/docker
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Security settings
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/mnt/rclone
ProtectHome=true

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target 